# kohya_ss-cu124/kohya_ss.service
[Unit]
Description=Kohya_ss GUI (CUDA 12.4 backend) — supports NVIDIA GPUs with Compute Capability ≥5.0 (Kepler/Maxwell/Pascal/Volta/Turing/Ampere/Ada)
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=15

[Service]
Type=simple
User=diffusion
Group=diffusion
WorkingDirectory=/opt/kohya_ss
ExecStart=/opt/kohya_ss/venv/bin/python /opt/kohya_ss/kohya_gui.py --listen 0.0.0.0 --server_port 7861 --headless
Restart=on-failure
RestartSec=5
TimeoutStartSec=0
TimeoutStopSec=30
KillSignal=SIGINT
KillMode=mixed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kohya_ss
UMask=007

[Install]
WantedBy=multi-user.target
