# comfyui/comfyui.service
[Unit]
Description=ComfyUI GUI (Latest PyTorch)
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=15

[Service]
Type=forking
User=diffusion
Group=diffusion
WorkingDirectory=/opt/ComfyUI
Environment=HOME=/opt/ComfyUI
Environment=BYOBU_PREFIX=/usr
Environment=TERM=xterm-256color
ExecStart=/usr/bin/byobu new-session -d -s comfyui 'bash -c "cd /opt/ComfyUI && /opt/ComfyUI/venv/bin/python main.py --listen 0.0.0.0 --port 8188"'
ExecStop=-/usr/bin/byobu kill-session -t comfyui
Restart=always
UMask=007

[Install]
WantedBy=multi-user.target
