# kohya_ss/kohya_ss.service
[Unit]
Description=Kohya_ss GUI (CUDA 12.8 backend) — ONLY for NVIDIA GPUs with Compute Capability ≥7.5 (Turing/Ampere/Ada/Hopper/Blackwell)
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=15

[Service]
Type=simple
User=diffusion
Group=diffusion
WorkingDirectory=/opt/kohya_ss
ExecStart=/opt/kohya_ss/venv/bin/python /opt/kohya_ss/kohya_gui.py --listen 0.0.0.0 --server_port 7861 --headless
Restart=on-failure
RestartSec=5
TimeoutStartSec=0
TimeoutStopSec=30
KillSignal=SIGINT
KillMode=mixed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kohya_ss
UMask=007

[Install]
WantedBy=multi-user.target
