# stable-diffusion-webui-forge-cu124/stable-diffusion-webui-forge.service
[Unit]
Description=Stable Diffusion WebUI Forge (backend: PyTorch 2.4 / CUDA 12.4)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=diffusion
Group=diffusion
WorkingDirectory=/opt/stable-diffusion-webui-forge

Environment=PYTHONUNBUFFERED=1

ExecStart=/opt/stable-diffusion-webui-forge/venv/bin/python /opt/stable-diffusion-webui-forge/launch.py --listen --no-half --port 7860

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=15

TimeoutStartSec=0
TimeoutStopSec=30
KillSignal=SIGINT
KillMode=mixed

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sdwebui-forge

UMask=007

[Install]
WantedBy=multi-user.target
