# stable-diffusion-webui-forge-cu124/stable-diffusion-webui-forge.service
[Unit]
Description=Stable Diffusion WebUI Forge (Experimental backend: PyTorch 2.4 / CU 12.4)
After=network.target

[Service]
Type=forking
User=diffusion
Group=diffusion
WorkingDirectory=/opt/stable-diffusion-webui-forge

Environment=PYTHONUNBUFFERED=1
Environment=BYOBU_PREFIX=/usr
Environment=TERM=xterm-256color

ExecStart=/usr/bin/byobu new-session -d -s sdwebui-forge 'bash -c "cd /opt/stable-diffusion-webui-forge && /opt/stable-diffusion-webui-forge/venv/bin/python launch.py --listen --no-half --disable-gpu-warning"'
ExecStop=-/usr/bin/byobu kill-session -t sdwebui-forge

Restart=always
StartLimitIntervalSec=120
StartLimitBurst=15
UMask=007

[Install]
WantedBy=multi-user.target
