# stable-diffusion-webui-forge-sd35/stable-diffusion-webui-forge-sd35.service
[Unit]
Description=Stable Diffusion WebUI Forge (sd35, python310)
After=network.target

[Service]
Type=simple
User=sdwebui-forge
Group=sdwebui-forge
WorkingDirectory=/opt/stable-diffusion-webui-forge-sd35

Environment=PYTHONUNBUFFERED=1

ExecStart=/opt/stable-diffusion-webui-forge-sd35/venv/bin/python /opt/stable-diffusion-webui-forge-sd35/launch.py --listen --no-half

Restart=always
StartLimitIntervalSec=120
StartLimitBurst=15
UMask=007

[Install]
WantedBy=multi-user.target