# stable-diffusion-webui-forge/stable-diffusion-webui-forge.service
[Unit]
Description=Stable Diffusion WebUI Forge (Stable backend: PyTorch 2.3.1 / CU 12.1)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=diffusion
Group=diffusion
WorkingDirectory=/opt/stable-diffusion-webui-forge

ExecStart=/opt/stable-diffusion-webui-forge/webui.sh

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=15

TimeoutStartSec=0
TimeoutStopSec=30
KillSignal=SIGINT
KillMode=mixed

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sdwebui-forge

UMask=007

[Install]
WantedBy=multi-user.target
